# systemd/user/rclone-documents.service
[Unit]
Description=OneDrive Documents Bisync
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.config/rclone/bisync-documents.sh
